# platforms/organelle_m/rootfs/usr/lib/systemd/system/cherrypy.service
[Unit] 
Description=Cherrypy Web Server

[Service] 
Type=simple 
WorkingDirectory=/home/music/fw_dir/web/server
ExecStart=/home/music/fw_dir/web/server/run.sh

[Install] 
WantedBy=multi-user.target
