# platforms/organelle_m/rootfs/usr/lib/systemd/system/createap.service
[Unit] 
Description=Create AP Daemon

[Service] 
Type=simple 
WorkingDirectory=/home/music/fw_dir/
ExecStart=/home/music/fw_dir/scripts/start-ap.sh

[Install] 
WantedBy=multi-user.target
