# platforms/organelle/rootfs/usr/lib/systemd/system/cherrypy.service
[Unit] 
Description=Cherrypy Web Server

[Service] 
Type=simple 
WorkingDirectory=/root/fw_dir/web/server
ExecStart=/root/fw_dir/web/server/run.sh

[Install] 
WantedBy=multi-user.target
