# platforms/organelle_m/rootfs/etc/systemd/system/ogweb.service
[Unit]
Description=Organelle Web
After=rc-local.service

[Service]
Type=simple
WorkingDirectory=/home/music/fw_dir/web
User=music
Group=music
ExecStart=/usr/bin/python -u /home/music/fw_dir/web/app.py 

# Restart on failure if needed
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target


